# bash/WEB-LOG-ANALYZER/attack-monitor.service
[Unit]
Description=Web Attack Monitor
After=network.target

[Service]
WorkingDirectory=/home/user/git_repo/Multi-Platform-Automation-Toolkit/bash/WEB-LOG-ANALYZER/
ExecStart=/home/user/git_repo/Multi-Platform-Automation-Toolkit/bash/WEB-LOG-ANALYZER/sentinel.sh
Restart=always
RestartSec=1200
User=root

[Install]
WantedBy=multi-user.target
